# bat-notif.service
[Service]
ExecStart=%h/.cargo/bin/bat-notif
Restart=always

[Install]
WantedBy=default.target
