# bat-notif.service
[Service]
ExecStart=/home/may/.cargo/bin/bat-notif
Restart=always

[Install]
WantedBy=default.target
